# dot-config/systemd/user/mpd-mpris.service
[Unit]
Description=Mako notifications
PartOf=graphical-session.target

[Service]
ExecStart=/usr/bin/mpd-mpris
Restart=on-failure

[Install]
WantedBy=sway-session.target