# dot-config/systemd/user/swayidle.service
[Unit]
Description=swayidle
PartOf=graphical-session.target

[Service]
ExecStart=/usr/bin/swayidle -w -C ~/.config/swayidle/config
Restart=on-failure

[Install]
WantedBy=sway-session.target