# dot-config/systemd/arch/mpd-default-defaut.service
[Unit]
Description=media player daemon

[Service]
Type=simple
ExecStart=mpd --systemd
ExecStop=mpd --kill

[Install]
WantedBy=default.target
